# templates/docker-cleanup.cron.j2
0 23 * * * root {{docker_cleanup.script_path}}/{{docker_cleanup.script_filename}} >> {{docker_cleanup.log}}
