# templates/docker-cleanup.cron.j2
0 23 * * * root {{docker_cleanup_script_file_path}} >> {{docker_cleanup_log}}
