# templates/docker-cleanup.timer.j2
[Unit]
Description=Run daily a docker cleanup

[Timer]
OnCalendar={{docker_cleanup_on_calendar}}
Persistent=true

[Install]
WantedBy=timers.target
